# extras/blackwood-update.service
[Unit]
Description=Update Blackwood from GitHub releases
Documentation=https://github.com/csweichel/blackwood/releases

[Service]
Type=oneshot
ExecStart=%h/.local/bin/blackwood-update
